# crontab
*  */3  * * *   root    /usr/bin/env bash -c 'cd /scripts/swbackup/ && source /scripts/swbackup/venv/bin/activate && /scripts/swbackup/venv/bin/python /scripts/swbackup/swbackup.py > /dev/null 2>&1'